# scripts/wg-pmtu-probe.service
[Unit]
Description=WireGuard Path MTU probe
After=network-online.target wg-quick@.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wg-pmtu-probe.sh
# A full probe of ~100 peers with tracepath fallbacks can take several
# minutes. 20 minutes is comfortable headroom; the script itself caps
# tracepath to 8s per peer so it can't runaway beyond ~100 × 16s ≈ 27m
# in the extreme worst case (all tracepath+ping misses).
TimeoutStartSec=1200
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

# Sandboxing: the script only needs to read WG state and write to
# /var/lib/wg-pmtu. Limit its access to the rest of the filesystem.
ProtectSystem=strict
ReadWritePaths=/var/lib/wg-pmtu
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
